# components/inseca-live-wks/etc/systemd/system/inseca-send-events.service
[Unit]
Description=Send INSECA events back home

[Service]
Type=oneshot
ExecStart=/opt/inseca-manager/send-events.py
