# components/inseca-live-wks/etc/systemd/system/inseca-send-events.timer
[Unit]
Description=Send INSECA events back home

[Timer]
OnUnitActiveSec=60s
OnBootSec=150s

[Install]
WantedBy=timers.target

